# deploy/content-publisher.service
[Unit]
Description=Content Publisher Platform
After=network-online.target postgresql.service
Wants=network-online.target

[Service]
Type=simple
User=content-publisher
Group=content-publisher
WorkingDirectory=/data/services/content-publisher
EnvironmentFile=/data/services/content-publisher/config/content-publisher.env
Environment=XDG_CONFIG_HOME=/data/services/content-publisher/data/config
ExecStart=/usr/bin/java -Xms256m -Xmx768m -Djava.io.tmpdir=/data/tmp/content-publisher -jar /data/services/content-publisher/app/content-publisher.jar
Restart=on-failure
RestartSec=5s
TimeoutStopSec=45s
SuccessExitStatus=143
UMask=0077
StandardOutput=append:/data/services/content-publisher/logs/application.log
StandardError=append:/data/services/content-publisher/logs/application.log
NoNewPrivileges=true
PrivateTmp=true
PrivateDevices=true
ProtectSystem=strict
ProtectHome=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictSUIDSGID=true
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
ReadWritePaths=/data/services/content-publisher/logs /data/services/content-publisher/data /data/tmp/content-publisher

[Install]
WantedBy=multi-user.target
